# roles/vps_edge_firewall/templates/vps-edge-firewall.service.j2
[Unit]
Description=VPS edge firewall nftables rules
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStartPre=-{{ edge_firewall_nft_bin }} delete table inet vps_edge_firewall
ExecStart={{ edge_firewall_nft_bin }} -f {{ edge_firewall_config_path }}
ExecStop=-{{ edge_firewall_nft_bin }} delete table inet vps_edge_firewall

[Install]
WantedBy=multi-user.target
