# roles/vps_authentik/templates/authentik-backup.service.j2
[Unit]
Description=Back up Edge Authentik
Requires=authentik-compose.service
After=authentik-compose.service

[Service]
Type=oneshot
ExecStart={{ authentik_base_dir }}/bin/backup.sh
