# roles/vps_authentik/templates/authentik-backup.timer.j2
[Unit]
Description=Run Edge Authentik backup

[Timer]
OnCalendar={{ authentik_backup_timer_on_calendar }}
Persistent=true

[Install]
WantedBy=timers.target
